# devops/systemd/deploy-caddy.path
[Unit]
Description="Monitor Caddyfile for changes"

[Path]
PathModified=/home/princetonpy/repos/princetonpy/caddy/Caddyfile
Unit=deploy-caddy.service

[Install]
WantedBy=multi-user.target